# root/etc/s6-overlay/s6-rc.d/init-taskchampion-sync-server-config/up
/etc/s6-overlay/s6-rc.d/init-taskchampion-sync-server-config/run
